# etc/sv/MyService1/run
#!/bin/zsh

# /usr/bin/fcitx5 &
# /home/fus/.fus/mount_external_ssd.sh

while sleep 6000; do :; done
